# incus/build-incus.timer
[Unit]
Description=Build Incus Container weekly

[Timer]
OnCalendar=weekly
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target