# tang/build-tang.timer
[Unit]
Description=Build Tang Container

[Timer]
OnCalendar=weekly
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target